# build/ansible/roles/jupyter/templates/upstart.conf.j2
description "ipython"

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 10 15

setuid {{ ipython_user }}
setgid {{ ipython_user }}

exec ipython notebook --no-browser
